# systemd_socket/99-usb-backup.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]", RUN+="/usr/local/bin/send_copy_trigger_signal.sh"