# systemd_socket/copy_trigger.socket
[Socket]
ListenStream=/tmp/copy_trigger.sock
SocketMode=0666